# docker-root/etc/s6-overlay/s6-rc.d/tigervnc/run
#!/command/with-contenv bash
sudo -Hu user /usr/bin/Xvnc :11 -geometry $VNC_GEOMETRY -depth 24
