# services.d/websocketify/run
#!/usr/bin/with-contenv bash
sudo -Hu user /app/src/websockify/run 9001 --wrap-mode=respawn -- /usr/bin/vncserver :1 -geometry $VNC_GEOMETRY -depth 24